# luci-app-maccloner/root/etc/hotplug.d/iface/99-maccloner
#!/bin/sh

[ "${ACTION}" = "ifup" ] || exit 0
[ "$(uci -q get maccloner.settings.enabled)" = "1" ] || exit 0

. /lib/functions/network.sh

DEVICE_NAME="${DEVICE:-$INTERFACE}"
LOCK_DIR="/var/run/maccloner"

logger -t maccloner "Hotplug triggered for $INTERFACE (device: $DEVICE_NAME)"

mkdir -p "$LOCK_DIR"

idx=0
while true; do
    rule_section="maccloner.@rule[${idx}]"
    uci -q get ${rule_section}.interface >/dev/null || break
    idx=$((idx + 1))

    [ "$(uci -q get ${rule_section}.enabled)" = "1" ] || continue

    RULE_IFACE=$(uci -q get ${rule_section}.interface)
    TARGET_SSID=$(uci -q get ${rule_section}.target_ssid)
    CUSTOM_MAC=$(uci -q get ${rule_section}.custom_mac)

    # Accept match on device name or logical interface name
    if [ "$DEVICE_NAME" != "$RULE_IFACE" ] && [ "$INTERFACE" != "$RULE_IFACE" ]; then
        continue
    fi

    [ -z "$CUSTOM_MAC" ] && continue

    # Lock to prevent re‑triggering within 20 seconds (enough for radio reload)
    LOCKFILE="$LOCK_DIR/$RULE_IFACE"
    if [ -f "$LOCKFILE" ] && [ $(($(date +%s) - $(stat -c %Y "$LOCKFILE"))) -lt 20 ]; then
        logger -t maccloner "Skipping $RULE_IFACE – recently modified"
        exit 0
    fi
    touch "$LOCKFILE"

    # Get current MAC (POSIX sed, no grep -P)
    CURRENT_MAC=$(ip link show "$RULE_IFACE" 2>/dev/null | sed -n 's/.*ether \([^ ]*\).*/\1/p' | tr '[:upper:]' '[:lower:]')
    if [ "$CURRENT_MAC" = "$(echo "$CUSTOM_MAC" | tr '[:upper:]' '[:lower:]')" ]; then
        logger -t maccloner "MAC already set on $RULE_IFACE, nothing to do."
        exit 0
    fi

    # Determine radio name (e.g., phy2‑sta0 → phy2 → radio2)
    PHY=""
    [ -x /usr/sbin/iw ] && PHY=$(iw dev "$RULE_IFACE" info 2>/dev/null | grep wiphy | awk '{print $2}')
    [ -z "$PHY" ] && PHY=$(cat /sys/class/net/$RULE_IFACE/phy80211/name 2>/dev/null)
    RADIO="radio${PHY#phy}"

    logger -t maccloner "Changing MAC of $RULE_IFACE to $CUSTOM_MAC"

    # Set MAC while interface is down
    ip link set dev "$RULE_IFACE" down
    ip link set dev "$RULE_IFACE" address "$CUSTOM_MAC"
    ip link set dev "$RULE_IFACE" up

    logger -t maccloner "MAC changed. Restarting radio $RADIO to reassociate..."

    # Reload only this radio – forces wpa_supplicant to re‑authenticate with new MAC
    wifi reload "$RADIO"

    logger -t maccloner "Radio $RADIO reloaded. DHCP will obtain IP."
    break
done
