# luci-app-maccloner/root/etc/hotplug.d/iface/99-maccloner
#!/bin/sh

[ "${ACTION}" = "ifup" ] || exit 0
[ "$(uci -q get maccloner.settings.enabled)" = "1" ] || exit 0

. /lib/functions/network.sh

logger -t maccloner "Hotplug triggered for $INTERFACE (device: $DEVICE)"

idx=0
while true; do
    rule_section="maccloner.@rule[${idx}]"
    uci -q get ${rule_section}.interface >/dev/null || break
    idx=$((idx + 1))

    ENABLED=$(uci -q get ${rule_section}.enabled)
    [ "$ENABLED" = "1" ] || continue

    IFACE=$(uci -q get ${rule_section}.interface)
    TARGET_SSID=$(uci -q get ${rule_section}.target_ssid)
    CUSTOM_MAC=$(uci -q get ${rule_section}.custom_mac)

    if [ "$INTERFACE" != "$IFACE" ] && [ "$DEVICE" != "$IFACE" ]; then
        continue
    fi

    [ -z "$CUSTOM_MAC" ] && { logger -t maccloner "No custom MAC, skipping."; continue; }

    # SSID verification
    if iwinfo "$IFACE" info >/dev/null 2>&1; then
        CURRENT_SSID=$(iwinfo "$IFACE" info | sed -n 's/.*SSID: "\(.*\)"/\1/p')
        if [ -n "$TARGET_SSID" ] && [ "$CURRENT_SSID" != "$TARGET_SSID" ]; then
            logger -t maccloner "SSID mismatch ($CURRENT_SSID != $TARGET_SSID), skipping."
            continue
        fi
    fi

    # Reliable PHY detection (works on all OpenWrt/ImmortalWrt)
    PHY=""
    if [ -x /usr/sbin/iw ]; then
        PHY=$(iw dev "$IFACE" info 2>/dev/null | grep wiphy | awk '{print $2}')
    fi
    if [ -z "$PHY" ] && [ -d /sys/class/net ]; then
        PHY=$(cat /sys/class/net/$IFACE/phy80211/name 2>/dev/null)
    fi
    RADIO="radio${PHY#phy}"

    logger -t maccloner "Applying custom MAC $CUSTOM_MAC to $IFACE (SSID: ${TARGET_SSID:-any}, radio: $RADIO)"

    # Loop over all wifi-iface sections using numeric indexing
    FOUND=0
    widx=0
    while true; do
        wifi_section="wireless.@wifi-iface[${widx}]"
        uci -q get ${wifi_section}.device >/dev/null || break
        widx=$((widx + 1))

        DEV=$(uci -q get ${wifi_section}.device)
        MODE=$(uci -q get ${wifi_section}.mode)
        SSID=$(uci -q get ${wifi_section}.ssid)

        if [ "$DEV" = "$RADIO" ] && [ "$MODE" = "sta" ] && [ "$SSID" = "$CURRENT_SSID" ]; then
            uci set ${wifi_section}.macaddr="$CUSTOM_MAC"
            FOUND=1
            break
        fi
    done

    if [ "$FOUND" -eq 1 ]; then
        uci commit wireless
        wifi reload
        logger -t maccloner "Wireless restarted for $IFACE. New MAC should be active."
    else
        logger -t maccloner "Could not find matching wifi-iface for $IFACE on $RADIO / SSID $CURRENT_SSID"
    fi
    break
done
