# luci-app-maccloner/files/etc/init.d/maccloner
#!/bin/sh /etc/rc.common

START=95
STOP=10

SERVICE_NAME=maccloner
PROG=/usr/sbin/wpa_cli
ACTION_SCRIPT=/etc/wpa_action.sh

start_service() {
    config_load wireless
    config_foreach start_wpa_cli wifi-iface
}

start_wpa_cli() {
    local mode device ifname disabled
    config_get mode "$1" mode
    config_get device "$1" device
    config_get disabled "$1" disabled

    [ "$disabled" = "1" ] && return
    [ "$mode" = "sta" ] || return

    # Determine the interface name; it's usually the radio's phyX-sta0
    local radio_phy
    radio_phy=$(uci get wireless.$device.phy 2>/dev/null || echo "${device#radio}")
    local sta_iface="${radio_phy}-sta0"

    # Only start if interface exists and wpa_supplicant control socket is present
    [ -S "/var/run/wpa_supplicant-${sta_iface}" ] || return

    procd_open_instance
    procd_set_param command "$PROG" -i "$sta_iface" -a "$ACTION_SCRIPT"
    procd_set_param respawn
    procd_close_instance
}

stop_service() {
    killall wpa_cli
}
